# dns/pihole-ftl/files/pihole-FTL.in
#!/bin/sh

# PROVIDE: pihole_FTL
#
# Add the following line to /etc/rc.conf to enable pihole_FTL:
#
# pihole_FTL_enable (bool): Set it to "YES" to enable pihole_FTL
#       Default is "NO".
#

. /etc/rc.subr

name="pihole_FTL"
rcvar="pihole_FTL_enable"

PIHOLE_USER="bind"
PIHOLE_GROUP="bind"

start_precmd="${name}_prestart"
start_cmd="${name}_start"
start_postcmd="echo done"
status_cmd="${name}_status"
stop_cmd="${name}_stop"

extra_commands="status"

load_rc_config $name 
: ${pihole_FTL_enable:=no}


pihole_FTL_prestart()
{
    # Touch files to ensure they exist (create if non-existing, preserve if existing)
    mkdir -p /var/run/pihole
    touch /var/run/pihole-FTL.pid /var/run/pihole-FTL.port
    chown $PIHOLE_USER:$PIHOLE_GROUP /var/run/pihole /var/run/pihole-FTL.pid /var/run/pihole-FTL.port

    mkdir -p /var/log/pihole
    chown $PIHOLE_USER:$PIHOLE_GROUP /var/log/pihole

    touch /usr/local/etc/pihole/dhcp.leases
    chown $PIHOLE_USER:$PIHOLE_GROUP /usr/local/etc/pihole/dhcp.leases

    /usr/local/bin/pihole-FTL --shm-cleanup
    rm -f /var/run/pihole/FTL.sock 2> /dev/null

    # Ensure that permissions are set so that pihole-FTL can edit all necessary files
    chown $PIHOLE_USER:$PIHOLE_GROUP /var/run/pihole-FTL.pid /var/run/pihole-FTL.port
    chown $PIHOLE_USER:$PIHOLE_GROUP /usr/local/etc/pihole /usr/local/etc/pihole/dhcp.leases 2> /dev/null
    chmod 0644 /var/log/pihole-FTL.log /var/run/pihole-FTL.pid /var/run/pihole-FTL.port /var/log/pihole.log
    
    # Chown database files to the user FTL runs as. We ignore errors as the files may not (yet) exist
    [ -f /usr/local/etc/pihole/pihole-FTL.db ] && chown $PIHOLE_USER:$PIHOLE_GROUP /usr/local/etc/pihole/pihole-FTL.db 2> /dev/null
    [ -f /usr/local/etc/pihole/gravity.db ] && chown $PIHOLE_USER:$PIHOLE_GROUP /usr/local/etc/pihole/gravity.db 2> /dev/null

    true
}
    
pihole_FTL_start()
{
    # TODO: how to port this in freebsd?
    #if setcap CAP_NET_BIND_SERVICE,CAP_NET_RAW,CAP_NET_ADMIN,CAP_SYS_NICE+eip "$(which pihole-FTL)"; then
    #  su -s /bin/sh -c "/usr/bin/pihole-FTL" "$FTLUSER"
    #else
    #  echo "Warning: Starting pihole-FTL as root because setting capabilities is not supported on this system"
    #  pihole-FTL
    #fi
    /usr/local/bin/pihole-FTL
}

is_running() {
    pgrep -o "pihole-FTL" > /dev/null
}

pihole_FTL_stop() {
  if is_running; then
    pkill -o pihole-FTL
    for i in {1..5}; do
      if ! is_running; then
        break
      fi

      echo -n "."
      sleep 1
    done
    echo

    if is_running; then
      echo "Not stopped; may still be shutting down or shutdown may have failed, killing now"
      pkill -o -9 pihole-FTL
      exit 1
    else
      echo "Stopped"
    fi
  else
    echo "Not running"
  fi
  echo
}

# Indicate the service status
pihole_FTL_status() {
  if is_running; then
    echo "[ ok ] pihole_FTL is running"
    exit 0
  else
    echo "[    ] pihole_FTL is not running"
    exit 1
  fi
}

run_rc_command "$1"
